# svc-monitor.timer
[Unit]
Description=Run SVC Monitor

[Timer]
OnCalendar=Hourly
RandomizedDelaySec=300

[Install]
WantedBy=timers.target

