# svc-monitor.service
[Unit]
Description=Run My SVC Monitor

[Service]
ExecStart=/home/jlindgren/git/tiny_projects/monitor/svc-monitor
Type=oneshot
